# 2.1.12-php5.6/scripts/supervisord-zmagento.conf
[program:magento2-setup]
command=sh -x /setup-magento.sh
startsecs=1
autorestart=false
startretries=1